# pi-puck-boot/debian/pi-puck-boot.prerm
#!/bin/sh -e
systemctl disable boot-hostname
#DEBHELPER#
